# systemd/drone-fusion.service
[Unit]
Description=Drone Fusion (audio + image -> GPIO)
# If the LED self-test unit is enabled, wait for it to finish before
# claiming the GPIO pins. Soft dependency: if the self-test unit is
# not installed, the main service still starts normally.
After=network-online.target sound.target drone-fusion-selftest.service
Wants=network-online.target

[Service]
Type=simple
# Adjust user/group/working dir to match your DietPi install.
User=dietpi
Group=dietpi
WorkingDirectory=/var/www/drone-fusion-pi
Environment=PYTHONUNBUFFERED=1
ExecStart=/var/www/drone-fusion-pi/venv/bin/python /var/www/drone-fusion-pi/main.py --no-audio
Restart=on-failure
RestartSec=5

# gpiozero with the lgpio backend uses /dev/gpiochip4 on Pi 5.
SupplementaryGroups=gpio audio video

[Install]
WantedBy=multi-user.target
